# dot_config/systemd/user/emacs.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/emacs emacsd --fg-daemon
