# dot_config/systemd/user/vdirsyncer.service.d/override.conf
[Unit]
BindsTo=davmail.service
After=davmail.service
